# launchd/io.homelab.storage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.homelab.storage</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>bash</string>
        <string>-lc</string>
        <string>sleep 30 &amp;&amp; echo "$(date): Starting storage mount service..." &amp;&amp; __HOME__/Documents/home-server/scripts/storage/setup_direct_mounts.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>USER</key>
        <string>__USER__</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/homelab-storage.out</string>
  <key>StandardErrorPath</key>
    <string>/tmp/homelab-storage.err</string>
</dict>
</plist>

