# launchd/io.homelab.updatecheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>io.homelab.updatecheck</string>
  <key>StartCalendarInterval</key><dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
  <key>ProgramArguments</key>
  <array><string>/bin/sh</string><string>-c</string><string>__REPO__/scripts/automation/check_updates.sh --apply</string></array>
  <key>StandardOutPath</key><string>/tmp/homelab-updatecheck.out</string>
  <key>StandardErrorPath</key><string>/tmp/homelab-updatecheck.err</string>
</dict></plist>
