# launchd/io.homelab.compose.immich.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>io.homelab.compose.immich</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /mnt/data/mac-mini-homeserver/services/immich && /opt/homebrew/bin/docker compose up -d</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/var/log/compose-immich.out.log</string>
  <key>StandardErrorPath</key><string>/var/log/compose-immich.err.log</string>
  <key>After</key><array><string>io.homelab.colima</string></array>
</dict>
</plist>
