# launchd/io.homelab.powermgmt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.homelab.powermgmt</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/Documents/home-server/scripts/ensure_power_settings.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StartInterval</key>
    <integer>3600</integer>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Volumes/warmstore/logs/powermgmt/powermgmt.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Volumes/warmstore/logs/powermgmt/powermgmt.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>USER</key>
        <string>__USER__</string>
    </dict>
</dict>
</plist>
